# Systemd Units/vk.service
[Unit]
Description=This Unit is needed to automatically start the Virtual Kubelet at system startup
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/interlink/bin
ExecStart=/etc/interlink/bin/vk
EnvironmentFile=/etc/interlink/.envs

[Install]
WantedBy=multi-user.target